# config/deployment/etc/services.d/beam/run
#!/usr/bin/with-contenv sh

echo "Starting MoodleNet backend..."

exec /opt/app/bin/${APP_NAME} start 
